# unibot.service
[Unit]
Description=unibot

[Service]
Type=simple
WorkingDirectory=/srv/unibot
ExecStart=/srv/unibot/unibot
Restart=always
RestartSec=90
StartLimitInterval=400
StartLimitBurst=3

[Unit]
After=default.target

[Install]
WantedBy=default.target
